# conf/systemd.service
[Unit]
Description=FlareSolverr
After=network.target

[Service]
Type=simple
User=__APP__
Group=__APP__
EnvironmentFile=__INSTALL_DIR__/.env
WorkingDirectory=__INSTALL_DIR__
ExecStart=__INSTALL_DIR__/flaresolverr

# Chromium needs /dev/shm for shared memory
PrivateTmp=yes
NoNewPrivileges=no

# Chromium sandbox requires these
ProtectSystem=full
ProtectHome=yes
ProtectClock=yes
ProtectHostname=yes
ProtectControlGroups=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
RestrictRealtime=yes
LockPersonality=yes
SystemCallArchitectures=native

# Chromium needs network and IPC access
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK

# Do NOT restrict namespaces or devices - Chromium needs them
# Do NOT add CAP_SYS_ADMIN to denied caps - Chromium sandbox needs it
# Do NOT filter syscalls too aggressively - Chromium needs clone, etc.

[Install]
WantedBy=multi-user.target
